# tomcat/supervisord.conf
[supervisord]
nodaemon = true

[unix_http_server]
file = /var/run/supervisor/supervisor.sock

[supervisorctl]
serverurl = unix:///var/run/supervisor/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:sshd]
command = /usr/sbin/sshd -D
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/supervisor/%(program_name)s.log

[program:tomcat]
command = /usr/local/bin/tomcat7-run
stopasgroup = true
redirect_stderr = true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
